# stage3/01-net-config/files/systemd/dnsmasq.service.d/after-online.conf
[Unit]
After=network-online.target
Wants=network-online.target
BindsTo=sys-subsystem-net-devices-wlan_ext.device
After=sys-subsystem-net-devices-wlan_ext.device