# OS-image/files/etc/systemd/system/hostapd.service.d/after-wlan-interface.conf
[Unit]
BindsTo=sys-subsystem-net-devices-wlan_ext.device
After=sys-subsystem-net-devices-wlan_ext.device
