# OS-image/files-lite/etc/systemd/user/ros.target
[Unit]
Description=ROS

[Install]
WantedBy=default.target
