# stage3/01-net-config/files/rc.local
#!/bin/sh -e

iptables-restore < /etc/iptables.ipv4.nat
exit 0